# tools/kalico-playground/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisord.log
pidfile=/run/supervisord.pid

# Our klippy + simulated MCUs, held alive for interactive use. runner.py is
# bind-mounted from the host (the serve-mode copy), so it can be iterated
# without rebuilding this image.
[program:klipper]
command=python3 /kalico/tools/kalico-sim/runner.py --serve --data-dir /printer_data
directory=/kalico
autostart=true
autorestart=true
startretries=999
stdout_logfile=/var/log/klipper.log
stderr_logfile=/var/log/klipper.log
priority=10

# Moonraker connects to /printer_data/klippy.sock (it retries until klippy is
# up, so start order is not critical).
[program:moonraker]
command=/opt/moonraker-env/bin/python /opt/moonraker/moonraker/moonraker.py -d /moonraker_data
directory=/opt/moonraker
autostart=true
autorestart=true
startretries=999
stdout_logfile=/var/log/moonraker.log
stderr_logfile=/var/log/moonraker.log
priority=20

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autostart=true
autorestart=true
stdout_logfile=/var/log/nginx-supervisor.log
stderr_logfile=/var/log/nginx-supervisor.log
priority=30
